# packages/audio/rpi-cirrus-config/udev.d/90-alsa-restore.rules
# SPDX-License-Identifier: GPL-2.0
# Copyright (C) 2016-present Team LibreELEC (https://libreelec.tv)

# Quando um dispositivo de som é detectado, restaura as configurações de volume
# Removido o parâmetro NAME que causava erro no udev moderno (Only network devices can be named)
SUBSYSTEM=="sound", KERNEL=="controlC*", ACTION=="add", GOTO="alsa_restore_go"
GOTO="alsa_restore_end"

LABEL="alsa_restore_go"

# Script de config separado para placa RPi-Cirrus
DRIVERS=="snd-rpi-cirrus", RUN+="rpi-cirrus-config", GOTO="alsa_restore_end"

# Script de config padrão para todas as outras placas
RUN+="soundconfig %k"

LABEL="alsa_restore_end"